# dashboard/openclaw-dashboard.service
[Unit]
Description=OpenClaw Trading Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pgre/.openclaw/dashboard
ExecStart=/usr/bin/node /home/pgre/.openclaw/dashboard/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
